# crontab
#@reboot /home/jsn/scripts/monitor_temps
#@reboot /home/jsn/scripts/check_ssh_tunnel
@reboot /home/jsn/scripts/start_fans
* * * * * /home/jsn/scripts/check_eth_proxy
